# overlay/etc/profile.d/bash_completion.sh
[[ $BASH && -f /usr/share/bash-completion/bash_completion ]] && \
  . /usr/share/bash-completion/bash_completion
