# overlay/etc/profile.d/bashrc.sh
if [ "$BASH" ]; then
  PS1="[\u@\h \W]\\$ "
fi
