# overlay/etc/profile.d/version.sh
#!/bin/sh

echo "Welcome to DockerRoot version $(cat /etc/version), $(docker -v)"
