# overlay/etc/profile.d/version.sh
#!/bin/sh

echo "DockerRoot version $(cat /etc/version), $(docker -v)"
